# deployment/redundancy-payments-service.upstart
description "Redundancy Payments Service"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 1 100

script
    /bin/bash <<EOT
sudo su ubuntu
source /etc/bash_completion.d/virtualenvwrapper
workon rps
redundancy_payments_service
EOT
end script